# plasma-kwin_wayland.service.in
[Unit]
Description=KDE Window Manager
PartOf=graphical-session.target

[Service]
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/kwin_wayland --xwayland
Type=notify
FileDescriptorStoreMax=200
Slice=session.slice
Restart=on-failure
