# plasma-kwin_wayland.service.in
[Unit]
Description=KDE Window Manager

[Service]
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/kwin_wayland
BusName=org.kde.KWin
KillMode=none
Slice=session.slice
