# plasma-kwin_x11.service.in
[Unit]
Description=KDE Window Manager
Wants=plasma-kcminit.service

[Service]
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/kwin_x11 --replace
BusName=org.kde.KWin
Slice=session.slice
